# .launchd/com.makerboarding.projectb_dw_db.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.makerboarding.projectb_dw_db</string>
      <key>ProgramArguments</key>
      <array>
        <string>/Users/dgdosen/.dotfiles/daily/cron_project_b_dw.sh</string>
      </array>
      <key>StartCalendarInterval</key>
      <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>14</integer>
      </dict>
    <key>StandardOutPath</key>
    <string>/Users/dgdosen/log/project_b_dw_db.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/dgdosen/log/project_b_dw_db.err</string>
    <key>SessionCreate</key>
    <true/>
  </dict>
</plist>
