# .launchd/com.agidevelopment.bear_seed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.agidevelopment.bear_seed</string>
      <!-- Monthly seed: creates daily/weekly/monthly/quarterly notes for a -->
      <!-- rolling ~6-month window. Runs on the 1st, just before bear_update. -->
      <key>ProgramArguments</key>
      <array>
        <string>/Users/dgdosen/.dotfiles/daily/bear_seed.sh</string>
      </array>
      <key>StartCalendarInterval</key>
      <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>5</integer>
      </dict>
    <key>StandardOutPath</key>
    <string>/Users/dgdosen/log/bear_seed.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/dgdosen/log/bear_seed.err</string>
    <key>SessionCreate</key>
    <true/>
  </dict>
</plist>
