# .launchd/org.danieldosen.update-do-dns.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>org.danieldosen.update-do-dns</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/dgdosen/.dotfiles/daily/update-do-dns.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/dgdosen</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/dgdosen</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/dgdosen/Library/Logs/update-dodns.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dgdosen/Library/Logs/update-dodns.err.log</string>
  </dict>
</plist>

